# crontab
#SHELL=/bin/bash
10 19 * * mon root /usr/local/bin/python3 /regionals.py  >> /cronlog 2>&1