# crontab
#SHELL=/bin/bash
31 19 * * mon root /usr/local/bin/python3 /webscraper.py  >> /cronlog 2>&1
